# system/scproxy.service.in
[Unit]
Description=Buypass SCProxy (unofficial)
Documentation=https://github.com/wvengen/scproxy
After=network.target

[Service]
ExecStart=$(BINDIR)/scproxy
User=$(RUNUSER)
PrivateUsers=true
